# deploy/ttc-daemon.service
[Unit]
Description=TTC Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ttc-automation
EnvironmentFile=/opt/ttc-automation/.env.server
ExecStart=/opt/ttc-automation/venv/bin/python ttc_daemon.py
Restart=always
RestartSec=5
User=ttc
Group=ttc

[Install]
WantedBy=multi-user.target
